# configs/arm-lab/system/home/ubuntu/.config/autostart/lxqt-desktop.desktop
[Desktop Entry]
Exec=pcmanfm-qt --desktop
Name=Desktop
Type=Application
Version=1.0
