# configs/arm-lab/system/home/ubuntu/.config/autostart/lxqt-policykit-agent.desktop
[Desktop Entry]
Exec=lxqt-policykit-agent
Name=PolicyKit Handler
Type=Application
Version=1.0
