# configs/riscv-lab/system/home/ubuntu/.config/autostart/lxqt-notifications.desktop
[Desktop Entry]
Exec=/usr/bin/lxqt-notificationd
Name=Notification Daemon
Type=Application
Version=1.0
