# configs/riscv-lab/system/home/ubuntu/.config/autostart/lxqt-panel.desktop
[Desktop Entry]
Exec=/usr/bin/lxqt-panel
Name=lxqt-panel
Type=Application
Version=1.0
