# configs/riscv-lab/system/home/ubuntu/.config/autostart/lxqt-openbox-session.desktop
[Desktop Entry]
Exec=/usr/bin/openbox-session
Name=openbox-session
Type=Application
Version=1.0
