# tools/system/usr/share/desktop/home/.config/autostart/lxqt-z-ibus.desktop
[Desktop Entry]
Exec=/usr/bin/ibus-daemon -d
Name=input-method
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
